# supervisor/schedule_work.ini
[program:schedule_work]
command=php artisan schedule:work
directory=/app/
autorestart=true
startsecs=3
startretries=10
stdout_logfile=/app/storage/logs/supervisor_log_work.out.log
stderr_logfile=/app/storage/logs/supervisor_log_work.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=www-data
priority=999
numprocs=1
process_name=%(program_name)s_%(process_num)02d